# k3s-portable-service/mkosi.extra/usr/lib/systemd/system/k3s.service
[Unit]
Description=Lightweight Kubernetes
Documentation=https://k3s.io
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=notify
EnvironmentFile=-/etc/rancher/k3s/env
KillMode=process
Delegate=yes
User=root
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
TimeoutStartSec=0
Restart=always
RestartSec=5s
ExecStartPre=-/sbin/modprobe br_netfilter
ExecStartPre=-/sbin/modprobe overlay
ExecStartPre=-/usr/local/bin/prepare.sh
StateDirectory=rancher kubelet cni
ConfigurationDirectory=rancher
LogsDirectory=containers pods
RuntimeDirectory=nri
BindReadOnlyPaths=/etc/passwd /etc/group
ExecStart=/usr/local/bin/k3s server $K3S_OPTS
